# tensorboard.conf
[program:tensorboard]
user=root
command=/opt/bin/start-tensorboard.sh 
process_name=%(program_name)s
stopsignal=INT
autostart=true
autorestart=true
redirect_stderr=true